# Software/RaspberryPi/AP_files/autostart/boot.desktop
[Desktop Entry]
Name=11ah_auto_start
Exec=lxterminal -e /home/pi/11ah_test/bc_boot.sh
Type=Application

